# ubuntu/18.04/custom/rc.local
#!/bin/bash

if lshw -c net | grep -q "driverversion=3.2.6-k"; then
    echo "Running old network driver - updating"

    MODDIR=$(ls -dSr /lib/modules/*/ | head -n 1)
    mkdir ${MODDIR}updates
    cp /opt/e1000e/e1000e.ko ${MODDIR}updates

    depmod
    update-initramfs -u

    if [ ! -f /var/log/e1000e.log ]; then
        shutdown -r +1 "Rebooting to update network driver"
    fi
    echo "$(date -Iseconds) Driver update to ${MODDIR}updates" >> /var/log/e1000e.log
fi

exit 0;
